# server/etc/systemd/system/multiroom_audio_rpi_server.service
[Unit]
Description=Multiroom Audio RPi Server Service

[Service]
Type=simple
ExecStart=snapserver

[Install]
WantedBy=multi-user.target
